# recipes-app/measy-iot/files/etc/systemd/system/connmand.service
[Unit]
Description=connman daemon

[Service]
Type=dbus
BusName=net.connman
ExecStart=/usr/bin/connmand -n


[Install]
WantedBy=multi-user.target
Alias=connmand.service

